# linux/etc/udev/rules.d/99-nvidia-hdmi-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card0", \
  RUN+="/bin/sh -c 'sleep 1 && echo detect > /sys/class/drm/card0/card0-HDMI-A-1/status'"
